# arrshole.service
[Unit]
Description=arrshole — stuck torrent monitor
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=smidge
WorkingDirectory=/home/smidge/development/arrshole
# Update PATH if your Node.js installation changes. Find it with: which node
Environment=PATH=/home/smidge/.nvm/versions/node/v25.5.0/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/usr/bin/env node dist/index.js
Restart=always
RestartSec=10
TimeoutStopSec=120
EnvironmentFile=/home/smidge/development/arrshole/.env
StandardOutput=journal
StandardError=journal
SyslogIdentifier=arrshole
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/smidge/development/arrshole
PrivateTmp=true

[Install]
WantedBy=multi-user.target
